# deploy/fra-bot.service
# systemd unit for the Fire & Rescue Academy Discord bot.
#
# Install (as the 'pi' user, adjust paths/user to taste):
#   sudo cp deploy/fra-bot.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now fra-bot
#
# Logs: journalctl -u fra-bot -f

[Unit]
Description=Fire & Rescue Academy Discord bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/FireAndRescueAcademyDiscordBot
ExecStart=/home/pi/FireAndRescueAcademyDiscordBot/.venv/bin/python -m fra_bot
Restart=on-failure
RestartSec=30

# Keep the Pi healthy: cap memory, lower priority slightly.
MemoryMax=512M
Nice=5

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/pi/FireAndRescueAcademyDiscordBot

[Install]
WantedBy=multi-user.target
